# courtauction-crawler.service
[Unit]
Description=Court Auction Crawler Service (FastAPI + Telegram Webhook)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/courtauction_crawler

# conda 환경 PATH 적용
Environment="PATH=/home/ubuntu/courtauction_crawler/env/bin:/usr/bin:/bin"
Environment="CONDA_PREFIX=/home/ubuntu/courtauction_crawler/env"
Environment="CONDA_DEFAULT_ENV=/home/ubuntu/courtauction_crawler/env"

ExecStart=/home/ubuntu/courtauction_crawler/env/bin/uvicorn main:app --host 0.0.0.0 --port 8000

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
